# supervisor/cron.d/laravel
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * cd /opt/www/smart_service && php artisan schedule:run >> /var/log/cron.log