# supervisor/conf.d/laravel.conf
[program:laravel]
directory=/opt/www/smart_service
command=php artisan queue:work --tries=1
autorestart=true

